# network-monitor.service
[Unit]
Description=Network Speed Test Monitor Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nivaldo
Group=nivaldo
WorkingDirectory=/home/nivaldo/src/projects/network-monitor
Environment=SPEEDTEST_VERBOSE=true
Environment=SPEEDTEST_INTERVAL=1800000
Environment=SPEEDTEST_MAX_RETRIES=3
Environment=SPEEDTEST_BACKOFF_DELAY=5000
Environment=SPEEDTEST_MAX_BACKOFF_DELAY=300000
Environment=SPEEDTEST_CIRCUIT_BREAKER_THRESHOLD=5
Environment=SPEEDTEST_CIRCUIT_BREAKER_TIMEOUT=1800000
ExecStart=/home/nivaldo/.bun/bin/bun run src/index.ts
Restart=always
RestartSec=10
StandardOutput=append:/var/log/network-monitor.log
StandardError=append:/var/log/network-monitor.error.log

[Install]
WantedBy=multi-user.target 